# estoraged@.service.in
[Unit]
Description=Encrypted storage daemon for %I

[Service]
Type=simple
Restart=always
ExecStart=@libexecdir@/eStoraged -b %I
