# service_files/mmc-bkops-enabled.target
[Unit]
Description=MMC BKOPs is enabled
