# service_files/xyz.openbmc_project.eStoraged.service.in
[Unit]
Description=Storage management daemon for an encrypted storage device
Wants=xyz.openbmc_project.EntityManager.service
After=xyz.openbmc_project.EntityManager.service

[Service]
Type=simple
Restart=always
ExecStart=@BINDIR@/eStoraged

[Install]
WantedBy=multi-user.target
